# templates/freeipa-backup.timer.j2
[Unit]
Description=Run daily backup of FreeIPA

[Timer]
{{ freeipa_backup_periodicity }}
Persistent=true

[Install]
WantedBy=timers.target
